# package/nodogsplash/files/nds.cron
# Submit nodogsplash statistics to control server
*/5 * * * * /etc/nodogsplash/push-statistics.sh
