# package/nodogsplash/files/nodogsplash.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=50

start() {
        nodogsplash -c /etc/nodogsplash/nodogsplash.conf -w /var/run/nodogsplash.sock
}

stop() {
        ndsctl -s /var/run/nodogsplash.sock stop
}
